# deploy/kis-backend.service
[Unit]
Description=KIS Django Backend
After=network.target

[Service]
User=admin
Group=admin
WorkingDirectory=/var/www/kis/backend
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"
EnvironmentFile=-/var/www/kis/backend/.env
ExecStart=/var/www/kis/backend/.venv/bin/gunicorn config.wsgi:application -c /var/www/kis/deploy/gunicorn.conf.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/www/kis /var/log/kis

[Install]
WantedBy=multi-user.target
