# deploy/kis-update.timer
[Unit]
Description=KIS Auto Update Timer - Check for updates every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
